# infra/scripts/paper-trading.service
[Unit]
Description=Paper Trading Bot
Documentation=https://github.com/hyperliquid-dex/hyperliquid-python-sdk
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu

# Working directory
WorkingDirectory=/home/ubuntu/hyperliquid-python-sdk

# Environment
Environment=PATH=/home/ubuntu/hyperliquid-python-sdk/venv/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/home/ubuntu/hyperliquid-python-sdk/.env

# Command to run
ExecStart=/home/ubuntu/hyperliquid-python-sdk/venv/bin/python -m scripts.paper_trading.scheduler

# Restart on failure
Restart=always
RestartSec=30

# Logging
StandardOutput=append:/home/ubuntu/hyperliquid-python-sdk/logs/paper_trading.log
StandardError=append:/home/ubuntu/hyperliquid-python-sdk/logs/paper_trading.log

# Resource limits
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
