# projects/dashboard/docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/app/logs

[program:indicator_manager]
command=python indicator_manager.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/app/logs/indicators.err.log
stdout_logfile=/app/logs/indicators.out.log
priority=1

[program:trigger_analyzer]
command=python triggers/analyzer.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/app/logs/trigger_analyzer.err.log
stdout_logfile=/app/logs/trigger_analyzer.out.log
priority=2

[program:trigger_streamer]
command=python triggers/streamer.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/app/logs/trigger_streamer.err.log
stdout_logfile=/app/logs/trigger_streamer.out.log
priority=3
startretries=10
startsecs=5

[program:paper_trader]
command=python paper_trading/paper_trader.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/app/logs/paper_trader.err.log
stdout_logfile=/app/logs/paper_trader.out.log
priority=4

[group:trading]
programs=indicator_manager,trigger_analyzer,trigger_streamer,paper_trader
priority=999